# backend/chess-backend.service
[Unit]
Description=CHESSbyOpus backend (FastAPI/uvicorn)
After=network.target

[Service]
Type=simple
User=chess
WorkingDirectory=/home/chess/chess-opus/backend
ExecStart=/home/chess/chess-opus/backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
